# deploy/freenote-backup.timer
[Unit]
Description=Daily Freenote backup (pg_dump + MinIO mirror)

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
