# deploy/freenote-backup.service
[Unit]
Description=Freenote backup (Postgres + MinIO)
After=network-online.target

[Service]
Type=oneshot
User=freenote
Group=freenote
ExecStart=/usr/local/bin/freenote-backup.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
StandardOutput=append:/var/log/freenote/backup.log
StandardError=append:/var/log/freenote/backup.log
